# cookbooks/encoders/templates/default/fileMonitorService.erb
#!/bin/bash
#
# chkconfig: 345 99 05
# description: Java deamon script


# SEE http://stackoverflow.com/questions/7687159/how-to-convert-a-java-program-to-daemon-with-jsvc
# SEE http://blog.widenhome.com/2011/06/28/how-to-install-jsvc-on-linux/

. /etc/init.d/functions

# Setup variables
EXEC=/usr/bin/jsvc
JAVA_HOME=/usr/lib/jvm/java
PROJECT_HOME=/data/apps/filemonitor
CLASS_PATH=$PROJECT_HOME/lib/commons-daemon-1.0.3.jar:$(echo $PROJECT_HOME/lib/*.jar | tr ' ' ':')
CLASS=com.clearchannel.customtalk.ingester.util.FileMonitor
USER=root
PIDFILE=/var/run/filemonitor.pid
CONFIG=-DconfigFile=/Utility/Config/filemonitors.cfg
LOG_OUT=/var/log/filemonitor/filemonitor.log
LOG_ERR=/var/log/filemonitor/filemonitor.error.log

# set umask perms
umask 0002

#echo $CLASS_PATH

do_exec()
{
	$EXEC   -home "$JAVA_HOME" -cp $CLASS_PATH $CONFIG  -user $USER -outfile $LOG_OUT -errfile $LOG_ERR -pidfile $PIDFILE $1 $CLASS
	ret=$?
	[ $ret -eq 0 ] && success || failure
    	echo
	chmod 644 ${PIDFILE} 
}

check_status() {
        if [ -r ${PIDFILE} ]; then
                PID=`cat $PIDFILE`
                PID_TEST=`ps -p $PID 2>&1 > /dev/null`
                if [ $? -ne 0 ]; then
                        echo "Stale PID file."
                else
                        echo -n " PIDFILE exists and is running ${CLASS} Under PID ${PID}"
                        success
                        echo
                fi
        else
                PID_APP=`pgrep -f "$CLASS" | tr '\n' ' '`
                if [ $? -eq 0 ]; then
                                echo -n " PIDFILE Does not exist but Running with PID ${PID_APP}"
                                failure
                                echo
                else
                        echo -n "$CLASS is not running."
                        failure
                        echo
                fi
        fi
}

case "$1" in
	start)
		echo -n "Starting Daemon: "
        	do_exec
            	;;
    	stop)
		echo -n "Stopping Daemon: "
        	do_exec "-stop"
            	;;
    	status)
        	check_status	
            	;;
    	restart)
		if [ -f "$PIDFILE" ]; then
			echo -n "Stopping Daemon: "
            		do_exec "-stop"
			echo -n "Starting Daemon: "
            		do_exec
		else
			echo -n "service not running, will do nothing: "
			failure
			echo
            		exit 1
		fi
            		;;
    	*)
        	echo "usage: daemon {start|stop|restart}" >&2
            	exit 3
            	;;
esac

